# systemd/yt-live.service
[Unit]
Description=YouTube Live
After=network.target

[Service]
ExecStart=/home/pi/.local/bin/yt-live
Restart=no

[Install]
WantedBy=default.target
