# systemd/yt-live.service
[Unit]
Description=YouTube Live
After=network.target

[Service]
User=pi
ExecStart=yt-live
Restart=no

[Install]
WantedBy=multi-user.target
